# app.service
[Unit]
Description = Diet server

[Service]
ExecStart = /var/www/app